# scripts/pariter-update.timer
[Unit]
Description=Daily Pariter update

[Timer]
OnCalendar=*-*-* 04:07:00
Persistent=true

[Install]
WantedBy=timers.target
